# build/crontab
# Rebuild all indexes daily.


*/5 * * * * php -f /var/www/html/magento2/bin/magento cron:run >> /var/www/html/magento2/var/log/magento.cron.log 2>&1;

*/30 * * * * php -f /var/www/html/magento2/bin/magento indexer:reindex >> /var/www/html/magento2/var/log/magento.indexer.log 2>&1;

#*/1 * * * * su magento -c "/usr/bin/php -c /etc/php/7.2/cli/php.ini /var/www/html/magento2/bin/magento cron:run > /var/www/html/magento2/var/log/magento.cron.log&""

#*/1 * * * * su magento -c "/usr/bin/php -c /etc/php/7.2/cli/php.ini /var/www/html/magento2/update/cron.php > /var/www/html/magento2/var/log/update.cron.log&""

#*/1 * * * * su magento -c "/usr/bin/php -c /etc/php/7.2/cli/php.ini /var/www/html/magento2/bin/magento setup:cron:run > /var/www/html/magento2/var/log/setup.cron.log&""


# An empty line is required at the end of this file for a valid cron file.


